# deploy/supervisor.conf
[unix_http_server]
file=/home/vagrant/supervisor.sock   ; (the path to the socket file)
chmod=0700                       ; sockef file mode (default 0700)

[supervisord]
logfile=/home/vagrant/supervisord.log ; (main log file;default $CWD/supervisord.log)
pidfile=/home/vagrant/supervisord.pid ; (supervisord pidfile;default supervisord.pid)
; childlogdir=/home/vagrant; ('AUTO' child log dir, default $TEMP)
user=vagrant

; the below section must remain in the config file for RPC
; (supervisorctl/web interface) to work, additional interfaces may be
; added by defining them in separate rpcinterface: sections
[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///home/vagrant/supervisor.sock ; use a unix:// URL  for a unix socket

; The [include] section can just contain the "files" setting.  This
; setting can list multiple files (separated by whitespace or
; newlines).  It can also contain wildcards.  The filenames are
; interpreted as relative to this file.  Included files *cannot*
; include files themselves.

[program:memex_explorer]
command = /vagrant/deploy/gunicorn_start.sh
user = vagrant
cwd = /vagrant/source
stdout_logfile = /home/vagrant/memex_explorer.log
redirect_stderr = true
autostart = true
killasgroup = true
