# scripts/supervisor_memex_explorer.conf
[unix_http_server]
file=/tmp/supervisor.sock                       ; path to your socket file

[supervisord]
logfile=./supervisord.log    ; supervisord log file
logfile_maxbytes=50MB                           ; maximum size of logfile before rotation
logfile_backups=10                              ; number of backed up logfiles
loglevel=error                                  ; info, debug, warn, trace
pidfile=./supervisord.pid                ; pidfile location
nodaemon=false                                  ; run supervisord as a daemon
minfds=1024                                     ; number of startup file descriptors
minprocs=200                                    ; number of process descriptors
user=vagrant                                       ; default user
childlogdir=./supervisord_child_logs               ; where child log files will live

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock         ; use a unix:// URL  for a unix socket

[program:memex_explorer]
command = ~/anaconda/envs/memex-explorer/bin/python  /vagrant/run.py
user = vagrant
stdout_logfile = ./supervisord_child_logs/memex_explorer.log
redirect_stderr = true
autostart = true
killasgroup = true

