# scripts/supervisor_memex_explorer.conf
[program:memex_explorer]
command = /opt/anaconda/bin/python /vagrant/run.py
user = vagrant
stdout_logfile = /var/log/memex_explorer.log
redirect_stderr = true
autostart = true
killasgroup = true
